# network-setup.service
[Unit]
Description=Configure network routes on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup_network.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
