# usbipd.service
[Unit]
Description=USBIP Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/usbip_start.sh
RemainAfterExit=yes
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

